# 5_busybox/lib/udev/rules.d/40-xdiagnose.rules
# do not edit this file, it will be overwritten on update

# Jesse Barnes on ubuntu-devel@lists.ubuntu.com:
#   You'll get three events, one when the error is detected, one before the
#   reset and one after.  Each has a different environment variable set; the
#   initial error has ERROR=1, the pre-reset event has RESET=1 and the
#   post-reset event has ERROR=0.

SUBSYSTEM=="drm", ACTION=="change", ENV{RESET}=="1", RUN+="/usr/share/apport/apport-gpu-error-intel.py"
